# .config/systemd/user/jupyterlab.service
[Unit]
Description=Jupyter Lab

[Service]
Type=simple
ExecStart=jupyter lab --no-browser

[Install]
WantedBy=default.target
